# dotfiles/.config/systemd/user/rsync-backup-snapshot.timer
[Unit]
After=network-online.target

[Timer]
OnBootSec=15min
OnUnitActiveSec=1day

[Install]
WantedBy=timers.target
